# deploy/systemd/basphere-api.service
[Unit]
Description=Basphere API Server
Documentation=https://github.com/basphere/basphere-api
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/basphere/basphere-api
ExecStart=/opt/basphere/basphere-api/build/basphere-api-linux-amd64 --config /etc/basphere/api.yaml
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=basphere-api

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
